# deploy/qar-analysis.service
[Unit]
Description=QAR Data Analysis System
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=exec
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/qar-analysis
Environment=PATH=/home/ubuntu/qar-analysis/venv/bin
Environment=PYTHONPATH=/home/ubuntu/qar-analysis/src
EnvironmentFile=/home/ubuntu/qar-analysis/config/production.env
ExecStart=/home/ubuntu/qar-analysis/venv/bin/uvicorn src.api.main:app --host 0.0.0.0 --port 8000 --workers 4
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qar-analysis

# 安全设置
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/ubuntu/qar-analysis/data /home/ubuntu/qar-analysis/logs /var/log/qar-analysis

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
